# BBB Scripts/Startup Script - Full Program/rc.local
#!/bin/sh -e

/usr/bin/X2_startCAN.sh
/home/debian/Sit_stand/canopend/./ALEX

exit 0
